# ansible/roles/temperature_sensor/files/supervisor.conf
[program:temperature_sensor]
command=/home/pi/rpi-mqtt-sensors/main.py
directory=/home/pi/rpi-mqtt-sensors
user=pi
autostart=true
autorestart=true
startsecs=5
