# ansible/roles/temperature_sensor/files/supervisor.conf
[program:temperature_sensor]
command=/home/pi/home/temperature/raspberrypi/main.py
directory=/home/pi/home/temperature/raspberrypi
user=pi
autostart=true
autorestart=true
